# redhat/rozofs-metadata-backup.timer
[Unit]
Description=Launch rozofs metadata backup on schedule

[Timer]
OnCalendar=*-*-* 02:00:00
Persistent=true

[Install]
WantedBy=timers.target
